# templates/systemd_unit.service.j2
Description={{ description }}
[Service]
ExecStart={{ exec_start_cmd }}
ExecStop={{ exec_kill_cmd }}
Type={{ service_type }}
[Install]
WantedBy=multi-user.target
